# debian/postinst
#!/bin/sh

set -e

# Create necessary directories
mkdir -p /etc/provisioner

# Set up systemd service
systemctl daemon-reload
systemctl enable provisioner.service
systemctl start provisioner.service

#DEBHELPER#