# debian/prerm
#!/bin/sh

set -e

# Stop and disable the service
systemctl stop provisioner.service
systemctl disable provisioner.service
systemctl daemon-reload

#DEBHELPER# 